# installer/core/ai-core.service
[Unit]
Description=RansomEye Ai Core
Documentation=https://ransomeye.v1.0/
Requires=correlation-engine.service
After=correlation-engine.service
BindsTo=correlation-engine.service
PartOf=ransomeye.target

[Service]
Type=notify
ExecCondition=/usr/bin/systemctl is-active correlation-engine.service --quiet
User=ransomeye
Group=ransomeye
WorkingDirectory=/opt/ransomeye
EnvironmentFile=-/opt/ransomeye/config/environment
Environment=RANSOMEYE_ORCHESTRATOR=systemd
ExecStart=/opt/ransomeye/bin/ransomeye-ai-core
Restart=on-failure
RestartSec=10
WatchdogSec=30
StartLimitInterval=300
StartLimitBurst=3
LimitNOFILE=65536
LimitNPROC=4096
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ransomeye/logs /opt/ransomeye/runtime
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ransomeye-ai-core

[Install]
WantedBy=ransomeye.target
